# archiso-radxa-x4/airootfs/etc/profile.d/plymouth-setup.sh
#!/bin/bash

# Set the Plymouth theme
plymouth-set-default-theme feralfile-splash

# Regenerate initramfs with Plymouth support
mkinitcpio -P

if [ -f /root/plymouth-setup.sh ]; then
  /root/plymouth-setup.sh
fi