# archiso-radxa-x4/airootfs/etc/systemd/system/boot-diagnostics.service
[Unit]
Description=Boot Graphics Diagnostics
After=plymouth-quit.service
Before=chromium-kiosk.service

[Service]
Type=oneshot
ExecStart=/root/boot-diagnostics.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target 