# archiso-radxa-x4/airootfs/etc/systemd/system/boot-setup.service
[Unit]
Description=FeralFile Boot Environment Setup
Before=multi-user.target
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/root/setup-boot-environment.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target 