# archiso-radxa-x4/airootfs/etc/systemd/system/chromium-kiosk.service
[Unit]
Description=Chromium Kiosk Mode

[Service]
Type=simple
User=feralfile
Group=feralfile
Environment="XDG_RUNTIME_DIR=/run/user/1000"
ExecStart=/home/feralfile/scripts/start-kiosk.sh
StandardOutput=append:/home/feralfile/.logs/chromium.log
StandardError=append:/home/feralfile/.logs/chromium.log
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target