# archiso-radxa-x4/airootfs/etc/systemd/system/connectd.service
[Unit]
Description=Feral File Connectd
After=local-fs.target

[Service]
Type=notify
ExecStart=/usr/bin/feral-connectd

[Install]
WantedBy=multi-user.target 