# archiso-radxa-x4/airootfs/etc/systemd/system/feral-watchdog.service
[Unit]
Description=FeralFile Kiosk Watchdog Daemon
After=chromium-kiosk.service network.target
Wants=network.target
ExecStartPre=/bin/sh -c 'while ! curl -s http://localhost:9222/json/version >/dev/null; do sleep 0.5; done'

[Service]
Type=notify
ExecStart=/usr/local/bin/feral-watchdog
Restart=always
RestartSec=5s
WatchdogSec=30s
User=root
Group=root

[Install]
WantedBy=multi-user.target