# archiso-radxa-x4/airootfs/etc/systemd/system/kiosk.service
[Unit]
Description=Chromium kiosk (Cage, 4K landscape)
After=network.target

[Service]
User=feralfile
PermissionsStartOnly=true
ExecStartPre=/bin/install -d -m0700 -o feralfile /run/user/1000
Environment=XDG_RUNTIME_DIR=/run/user/1000
SupplementaryGroups=video input seat
ExecStart=/home/feralfile/bin/kiosk-launch.sh
Restart=always

[Install]
WantedBy=multi-user.target
