# archiso-radxa-x4/airootfs/etc/systemd/system/seatd.service.d/override.conf
[Service]
ExecStartPost=/bin/chmod 660 /run/seatd.sock
ExecStartPost=/bin/chgrp seat /run/seatd.sock 