# archiso-radxa-x4/airootfs/etc/skel/.config/autostart/feralfile-launcher.desktop
[Desktop Entry]
Type=Application
Name=FeralFile Launcher
Exec=/usr/bin/feralfile-launcher
Terminal=false
X-GNOME-Autostart-enabled=true
