# configs/99-display-hotplug.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/usr/local/bin/dynamic_resize.sh"